# config/.config/autostart/bspwm.desktop
[Desktop Entry]
Name=BSPWM
Comment=Tiling Window Manager
Exec=bspwm
Type=Application
Terminal=false
StartupNotify=false
